# webservice/.add/etc/profile.d/30-local-pgbouncer.sh

# Сюда прописываются настройки доступа к БД через pgbouncer
# реальные настройки подключения сохранены в /etc/pgbouncer/pgbouncer.ini

